# etc/crontab/macholic@digital-bot.crontab
* * * * * /home/macholic/project/InfoWSD/bin/pcroom_usage.pl
